# scripts/systemd/brain-sync.service
[Unit]
Description=Sync $ACA_DATA git repository
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.polecat/aops/scripts/brain-sync.sh
# Use the installed location in ~/.polecat/aops (not worktree)
